# wg-forge-checklimits.service
[Unit]
Description=wg-forge bandwidth limit enforcement

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wg-forge checklimits
